# systemd/.config/systemd/user/backup.timer
[Unit]
Description=Periodic backups

[Timer]
OnCalendar=18:00
Persistent=true

[Install]
WantedBy=timers.target