# systemd/.config/systemd/user/backup.timer
[Unit]
Description=Periodic backups

[Timer]
OnCalendar=Mon,Thu *-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
