# systemd/.config/systemd/user/backup-nas.timer
[Unit]
Description=Periodic backups

[Timer]
OnCalendar=15:00
Persistent=true

[Install]
WantedBy=timers.target