# systemd/.config/systemd/user/backup-usbkey.timer
[Unit]
Description=Periodic USB backups

[Timer]
OnCalendar=11:00
Persistent=true

[Install]
WantedBy=timers.target
