# systemd/.config/systemd/user/move_screenshots.timer
[Unit]
Description=Periodically move screenshots

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
