# systemd/.config/systemd/user/move_screenshots.service
[Unit]
Description=Move screenshots out of the way

[Service]
ExecStart=
ExecStart="/bin/sh" "-c" "mv -f $HOME/Pictures/Screenshot* /tmp"