# systemd/.config/systemd/user/update-zsh-autosuggestions.timer
[Unit]
Description=Periodically update zsh-autosuggestions

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
