# systemd/.config/systemd/user/update_protonvpn.timer
[Unit]
Description=Periodically update ProtonVPN

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
