# systemd/.config/systemd/user/update_fzf.service
[Unit]
Description=Update fzf

[Service]
ExecStart=
ExecStart="/bin/sh" "-c" "cd $HOME/.fzf && git pull && ./install --bin"