# systemd/.config/systemd/user/update_fzf.timer
[Unit]
Description=Periodically update fzf

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
