# systemd/.config/systemd/user/update-fzf-tab.timer
[Unit]
Description=Periodically update fzf-tab

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
