# systemd/.config/systemd/user/update-fzf-tab.service
[Unit]
Description=Update fzf-tab

[Service]
ExecStart=
ExecStart="/bin/sh" "-c" "cd $HOME/src/github.com/Aloxaf/fzf-tab && git pull"